# data/supervisord.conf
[supervisord]
nodaemon=true
; Don't suppress script output, but reduce supervisor verbosity
loglevel=debug

[program:Updater]
; Use shell directly to bypass permission issues
command=sh -c "sh /data/entrypoint.sh | awk '/./ {print \"[     UPDATER     ]\", $0}'"
autorestart=true

[program:api-dashboard]
; Same pattern as Updater
command=sh -c "pwsh -f /data/www/Start-Server.ps1 | awk '/./ {print \"[ API / Dashboard ]\", $0}'"
autorestart=true